# sysutils/amazon-ssm-agent/pkg-plist
bin/ssm-agent-worker
etc/rc.d/amazon-ssm-agent
etc/amazon/ssm/seelog_unix.xml.template
etc/amazon/ssm/amazon-ssm-agent.json.template
sbin/ssm-document-worker
sbin/amazon-ssm-agent
sbin/ssm-session-logger
sbin/ssm-cli
sbin/ssm-session-worker
@dir /var/run/amazon/ssm
@dir /var/run/amazon
@dir /var/log/amazon
@dir /var/lib/amazon/ssm/daemons
@dir /var/lib/amazon/ssm/localcommands
@dir /var/lib/amazon/ssm/ipc
@dir /var/lib/amazon/ssm/locks
@dir /var/lib/amazon/ssm/manifests
@dir /var/lib/amazon/ssm/packages
@dir /var/lib/amazon/ssm
@dir /var/lib/amazon
@dir etc/amazon/ssm
@dir etc/amazon
