# branches/2018Q4/comms/comserv/pkg-plist
man/man8/comservd.8.gz
sbin/comservd
etc/comservd.conf.sample
etc/rc.d/comservd.sh.sample
