# tags/RELEASE_0_3_0/net/mDNSResponder/files/mdnsd.sh.in
#!/bin/sh
# $FreeBSD: ports/net/mDNSResponder/files/mdnsd.sh.in,v 1.6 2010/03/27 00:13:47 dougb Exp $
# $MidnightBSD$

# PROVIDE: mdnsd
# REQUIRE: DAEMON
# KEYWORD: shutdown

mdnsd_enable="${mdnsd_enable-NO}"

. /etc/rc.subr

name=mdnsd
rcvar=`set_rcvar`
command="%%PREFIX%%/sbin/${name}"
pidfile="/var/run/${name}.pid"

load_rc_config $name

run_rc_command $*
