# sysutils/midnightbsd-stats/pkg-plist
etc/periodic/monthly/310.statistics
etc/rc.d/midnightbsd-stats.sh
@dirrmtry etc/periodic/monthly
@dirrmtry etc/periodic
