# trunk/mail/dkim-milter/pkg-plist
%%STATS%%$MidnightBSD: mports/Mk/components/maintainer.mk,v 1.2 2009/03/26 22:19:02 ctriv Exp $
bin/dkim-genkey
bin/dkim-testkey
bin/dkim-testssp
etc/rc.d/milter-dkim.sh
etc/mail/dkim-filter.conf.sample
include/dkim.h
libexec/dkim-filter
%%DOCSDIR%%/FEATURES
%%DOCSDIR%%/INSTALL
%%DOCSDIR%%/KNOWNBUGS
%%DOCSDIR%%/LICENSE
%%DOCSDIR%%/README
%%DOCSDIR%%/RELEASE_NOTES
%%DOCSDIR%%/autorespond.csh
%%DOCSDIR%%/draft-ietf-dkim-deployment-04.txt
%%DOCSDIR%%/draft-ietf-dkim-overview-12.txt
%%DOCSDIR%%/draft-ietf-dkim-rfc4871-errata-05.txt
%%DOCSDIR%%/draft-ietf-dkim-ssp-10.txt
%%DOCSDIR%%/draft-kucherawy-dkim-reporting-04.txt
%%DOCSDIR%%/draft-shafranovich-feedback-report-07.txt
%%DOCSDIR%%/rfc4686.txt
%%DOCSDIR%%/rfc4871.txt
%%DOCSDIR%%/rfc5016.txt
%%DOCSDIR%%/rfc5451.txt
%%DOCSDIR%%/rfc5518.txt
%%DOCSDIR%%/README.dkim-filter
%%DOCSDIR%%/README.libdkim
%%DOCSDIR%%/libdkim/overview.html
%%DOCSDIR%%/libdkim/dkim.html
%%DOCSDIR%%/libdkim/dkim_alg_t.html
%%DOCSDIR%%/libdkim/dkim_body.html
%%DOCSDIR%%/libdkim/dkim_canon_t.html
%%DOCSDIR%%/libdkim/dkim_cbstat.html
%%DOCSDIR%%/libdkim/dkim_chunk.html
%%DOCSDIR%%/libdkim/dkim_close.html
%%DOCSDIR%%/libdkim/dkim_dnssec.html
%%DOCSDIR%%/libdkim/dkim_eoh.html
%%DOCSDIR%%/libdkim/dkim_eom.html
%%DOCSDIR%%/libdkim/dkim_flush_cache.html
%%DOCSDIR%%/libdkim/dkim_free.html
%%DOCSDIR%%/libdkim/dkim_get_user_context.html
%%DOCSDIR%%/libdkim/dkim_getcachestats.html
%%DOCSDIR%%/libdkim/dkim_getdomain.html
%%DOCSDIR%%/libdkim/dkim_geterror.html
%%DOCSDIR%%/libdkim/dkim_policy_getdnssec.html
%%DOCSDIR%%/libdkim/dkim_getmode.html
%%DOCSDIR%%/libdkim/dkim_getpolicystr.html
%%DOCSDIR%%/libdkim/dkim_getpresult.html
%%DOCSDIR%%/libdkim/dkim_getpresultstr.html
%%DOCSDIR%%/libdkim/dkim_getresultstr.html
%%DOCSDIR%%/libdkim/dkim_getsighdr.html
%%DOCSDIR%%/libdkim/dkim_getsiglist.html
%%DOCSDIR%%/libdkim/dkim_getsignature.html
%%DOCSDIR%%/libdkim/dkim_header.html
%%DOCSDIR%%/libdkim/dkim_init.html
%%DOCSDIR%%/libdkim/dkim_key_syntax.html
%%DOCSDIR%%/libdkim/dkim_lib.html
%%DOCSDIR%%/libdkim/dkim_minbody.html
%%DOCSDIR%%/libdkim/dkim_ohdrs.html
%%DOCSDIR%%/libdkim/dkim_options.html
%%DOCSDIR%%/libdkim/dkim_param_t.html
%%DOCSDIR%%/libdkim/dkim_policy.html
%%DOCSDIR%%/libdkim/index.html
%%DOCSDIR%%/libdkim/dkim_policy_getreportinfo.html
%%DOCSDIR%%/libdkim/dkim_policy_syntax.html
%%DOCSDIR%%/libdkim/dkim_policy_t.html
%%DOCSDIR%%/libdkim/dkim_presult.html
%%DOCSDIR%%/libdkim/dkim_pstate.html
%%DOCSDIR%%/libdkim/dkim_query_t.html
%%DOCSDIR%%/libdkim/dkim_set_dns_callback.html
%%DOCSDIR%%/libdkim/dkim_set_final.html
%%DOCSDIR%%/libdkim/dkim_set_key_lookup.html
%%DOCSDIR%%/libdkim/dkim_set_margin.html
%%DOCSDIR%%/libdkim/dkim_set_policy_lookup.html
%%DOCSDIR%%/libdkim/dkim_set_prescreen.html
%%DOCSDIR%%/libdkim/dkim_set_signature_handle.html
%%DOCSDIR%%/libdkim/dkim_set_signature_handle_free.html
%%DOCSDIR%%/libdkim/dkim_set_signature_tagvalues.html
%%DOCSDIR%%/libdkim/dkim_set_signer.html
%%DOCSDIR%%/libdkim/dkim_set_user_context.html
%%DOCSDIR%%/libdkim/dkim_sig_getbh.html
%%DOCSDIR%%/libdkim/dkim_sig_getcanonlen.html
%%DOCSDIR%%/libdkim/dkim_sig_getcontext.html
%%DOCSDIR%%/libdkim/dkim_sig_getdnssec.html
%%DOCSDIR%%/libdkim/dkim_sig_getdomain.html
%%DOCSDIR%%/libdkim/dkim_sig_geterror.html
%%DOCSDIR%%/libdkim/dkim_sig_geterrorstr.html
%%DOCSDIR%%/libdkim/dkim_sig_getflags.html
%%DOCSDIR%%/libdkim/dkim_sig_getkeysize.html
%%DOCSDIR%%/libdkim/dkim_sig_getreportinfo.html
%%DOCSDIR%%/libdkim/dkim_sig_getselector.html
%%DOCSDIR%%/libdkim/dkim_sig_getsignalg.html
%%DOCSDIR%%/libdkim/dkim_sigerror.html
%%DOCSDIR%%/libdkim/dkim_sig_getsigntime.html
%%DOCSDIR%%/libdkim/dkim_sig_hdrsigned.html
%%DOCSDIR%%/libdkim/dkim_sig_ignore.html
%%DOCSDIR%%/libdkim/dkim_sig_process.html
%%DOCSDIR%%/libdkim/dkim_sig_syntax.html
%%DOCSDIR%%/libdkim/dkim_siginfo.html
%%DOCSDIR%%/libdkim/dkim_sigkey_t.html
%%DOCSDIR%%/libdkim/dkim_sign.html
%%DOCSDIR%%/libdkim/dkim_ssl_version.html
%%DOCSDIR%%/libdkim/dkim_stat.html
%%DOCSDIR%%/libdkim/dkim_verify.html
%%DOCSDIR%%/libdkim/rfc2822_mailbox_split.html
@dirrm etc/mail
@dirrm %%DOCSDIR%%/libdkim
@dirrm %%DOCSDIR%%
