# tags/1.1/emulators/hyperv-is/pkg-plist
hyperv/scripts/hv_get_dhcp_info
hyperv/scripts/hv_get_dns_info
hyperv/scripts/hv_set_ifconfig
@cwd /
boot/modules/hv_ata_pci_disengage.ko
boot/modules/hv_netvsc.ko
boot/modules/hv_storvsc.ko
boot/modules/hv_utils.ko
boot/modules/hv_vmbus.ko
@cwd %%LOCALBASE%%
etc/rc.d/hv_kvpd
man/man1/hyperv-is.1.gz
man/man4/hv_ata_pci_disengage.4.gz
man/man4/hv_kvp.4.gz
man/man4/hv_netvsc.4.gz
man/man4/hv_storvsc.4.gz
man/man4/hv_utils.4.gz
man/man4/hv_vmbus.4.gz
